# rknpu2/runtime/Linux/package/DEBIAN/postinst
#!/bin/sh
set -e

if [ "$1" = "configure" ]; then
    ldconfig || true
    if command -v systemctl >/dev/null 2>&1; then
        systemctl daemon-reload || true
        systemctl enable rknn-server.service || true
        systemctl start rknn-server.service || true
    fi
fi

exit 0
